# worklog-timer.service
[Unit]
Description=Worklog Timer - Periodic work activity logger
# Bind to the graphical session so the daemon starts only once the
# desktop is up and stops with it. DISPLAY/XAUTHORITY are re-discovered
# at runtime before every popup, so no Environment= lines are needed.
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=@WORKLOG_PATH@/worklog start --foreground --interval 45 --anchor 08:45
Restart=on-failure
RestartSec=10

[Install]
WantedBy=graphical-session.target
