# worklog-timer.service
[Unit]
Description=Worklog Timer - Periodic work activity logger
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=@WORKLOG_PATH@/worklog start --foreground --interval 45
Restart=on-failure
RestartSec=10
# Pass graphical session env vars if available; daemon discovers XAUTHORITY at runtime
Environment=DISPLAY=:0

[Install]
WantedBy=default.target
